# templates/default/cfssl.service.erb
[Unit]
Description=CloudFlare SSL Server Daemon
After=network-online.target

[Service]
Type=simple
User=<%= node['cfssl']['client']['user'] %>
Group=<%= node['cfssl']['client']['group'] %>
UMask=007
ExecStart=<%= node['cfssl']['install_path'] %>/cfssl serve -address <%= node['cfssl']['client']['bind'] %> -port <%= node['cfssl']['client']['port'] %> -remote <%= node['cfssl']['remote']['address'] %>:<%= node['cfssl']['remote']['port'] %> -config <%= node['cfssl']['config_path'] %>/conf/cfssl.json

Restart=on-failure

TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
